# deploy/systemd/jagc-restart.service
[Unit]
Description=Restart jagc (triggered by jagc.path)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart jagc.service
