# deploy/systemd/jagc.service
[Unit]
Description=jagc server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jagc
Group=jagc
WorkingDirectory=/var/lib/jagc/workspace
EnvironmentFile=/etc/jagc/jagc.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env node /usr/local/lib/node_modules/jagc/dist/server/main.mjs
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2
TimeoutStopSec=30
KillSignal=SIGINT

# Logs end up in journald. View with:
#   journalctl -u jagc -f

[Install]
WantedBy=multi-user.target
