# deploy/systemd/kaiser-restart.service
[Unit]
Description=Restart kaiser (triggered by kaiser.path)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart kaiser.service
