# deploy/launchd/com.jagc.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jagc.server</string>

    <!--
      Prefer `jagc install` for managed setup.
      This file is the manual fallback when operators need to hand-edit launchd config.
    -->
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/node</string>
      <string>/usr/local/lib/node_modules/jagc/dist/server/main.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/you/.jagc</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>JAGC_WORKSPACE_DIR</key>
      <string>/Users/you/.jagc</string>
      <key>JAGC_DATABASE_PATH</key>
      <string>/Users/you/.jagc/jagc.sqlite</string>
      <key>JAGC_HOST</key>
      <string>127.0.0.1</string>
      <key>JAGC_PORT</key>
      <string>31415</string>
      <key>JAGC_RUNNER</key>
      <string>pi</string>
      <key>JAGC_LOG_LEVEL</key>
      <string>info</string>
      <!-- Optional provider env vars -->
      <!-- <key>OPENAI_API_KEY</key><string>...</string> -->
      <!-- Optional Telegram polling -->
      <!-- <key>JAGC_TELEGRAM_BOT_TOKEN</key><string>...</string> -->
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/you/.jagc/logs/server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/you/.jagc/logs/server.err.log</string>
  </dict>
</plist>
